# roles/laptop/files/home-gmas-nas.mount
[Unit]
Description = NAS

[Mount]
What = nas:/nas
Where = /home/gmas/nas
Type = nfs4

[Install]
WantedBy = multi-user.target

